# scripts/systemd/copilot-workspace-docker.service.example
# ~/.config/systemd/user/copilot-workspace-docker.service にコピーして使う
[Unit]
Description=Copilot workspace container
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=%h/development/copilot-workspace-docker
ExecStart=/usr/bin/docker compose up -d workspace
ExecStop=/usr/bin/docker compose down
TimeoutStartSec=0

[Install]
WantedBy=default.target
